# .config/autostart/nemo.desktop
[Desktop Entry]
Name=Nemo
Exec=nemo
Terminal=false
Type=Application
StartupNotify=true

